# upstart_configs/serverphp_runner.conf
description "runs server.php"

start on startup
stop on shutdown
respawn

script
 php /var/www/html/demovpn-new/app/socket/server.php >> server.php.log.out 2>> server.php.log.err 
end script
